# service/i3.service
[Unit]
Description=i3 Window Manager
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/i3

[Install]
RequiredBy=xsession.target
