# service/jupyter.socket
[Unit] 
Description = broken

[Socket] 
ListenStream = 4444
Accept = yes

[Install] 
WantedBy = sockets.target
