# service/mirrorlist.service
[Unit]
Description=Update pacman mirrorlist

[Service]
Type=oneshot
ExecStart=/usr/var/local/tools/root/util/update_pacman_mirrorlist.sh
