# service/pushrsync.service
[Unit]
Description=Rsync everything

[Service]
Type=oneshot
ExecStart=/usr/var/local/tools/root/bin/pushrsync
